# src/xrt/targets/broker/broker.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>@XRT_IPC_MACH_SERVICE_NAME@</string>
	<key>MachServices</key>
	<dict>
		<key>@XRT_IPC_MACH_SERVICE_NAME@</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>$<TARGET_FILE:monado-broker></string>
	</array>
	<key>StandardOutPath</key>
	<string>/var/tmp/@XRT_IPC_MACH_SERVICE_NAME@.log</string>
	<key>StandardErrorPath</key>
	<string>/var/tmp/@XRT_IPC_MACH_SERVICE_NAME@-error.log</string>
</dict>
</plist>
